# systemd/system/fstrim.service
[Unit]
Description=Discard unused blocks on filesystems from time to time

[Service]
Type=oneshot
ExecStart=/usr/sbin/fstrim --all
