# systemd/system/security-updates.service
[Unit]
Description=Automated Security Updates

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Syu --noconfirm --needed
