# systemd/system/security-updates.timer
[Unit]
Description=Run Security Updates Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
